# deploy/com.neve.server.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd 服务模板:@占位符@ 由 `make install-service` 渲染,不要直接拷贝本文件 -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.neve.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>@NEVE_ROOT@/neve</string>
    </array>

    <key>WorkingDirectory</key>
    <string>@NEVE_ROOT@</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NEVE_DATA_DIR</key>
        <string>@NEVE_ROOT@/data</string>
        <key>NEVE_PORT</key>
        <string>9999</string>
        <!-- 显式钉死时区,避免系统时区随旅行自动切换导致日期归属漂移 -->
        <key>TZ</key>
        <string>Asia/Singapore</string>
        <!-- 无感记账入口:值来自 deploy/local.env(不入库),留空则 /api/inbox 关闭 -->
        <key>NEVE_INBOX_TOKEN</key>
        <string>@NEVE_INBOX_TOKEN@</string>
        <key>NEVE_AI_PROVIDER</key>
        <string>@NEVE_AI_PROVIDER@</string>
        <key>NEVE_AI_API_KEY</key>
        <string>@NEVE_AI_API_KEY@</string>
        <key>NEVE_AI_MODEL</key>
        <string>@NEVE_AI_MODEL@</string>
        <key>NEVE_BARK_URL</key>
        <string>@NEVE_BARK_URL@</string>
        <!-- 数据备份:填了远程 URL 才启用。服务端把账本镜像进 iCloud 外的 git 仓库并推送,
             绕开 launchd 沙箱对快捷指令 iCloud 容器的 readdir/chdir 限制。NEVE_BACKUP_DIR
             留空用默认 ~/Library/Application Support/Neve/data-backup -->
        <key>NEVE_BACKUP_REMOTE</key>
        <string>@NEVE_BACKUP_REMOTE@</string>
        <key>NEVE_BACKUP_DIR</key>
        <string>@NEVE_BACKUP_DIR@</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <!-- 崩溃后自动重启,正常退出不拉起 -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- 后台低优先级,不抢占前台资源 -->
    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>1024</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>@HOME@/Library/Logs/neve.log</string>

    <key>StandardErrorPath</key>
    <string>@HOME@/Library/Logs/neve.error.log</string>
</dict>
</plist>
