# deploy/com.neve.server.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd 服务模板:@占位符@ 由 `make install-service` 渲染,不要直接拷贝本文件 -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.neve.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>@NEVE_ROOT@/neve</string>
    </array>

    <key>WorkingDirectory</key>
    <string>@NEVE_ROOT@</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NEVE_DATA_DIR</key>
        <string>@NEVE_ROOT@/data</string>
        <key>NEVE_PORT</key>
        <string>9999</string>
        <!-- 显式钉死时区,避免系统时区随旅行自动切换导致日期归属漂移 -->
        <key>TZ</key>
        <string>Asia/Singapore</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <!-- 崩溃后自动重启,正常退出不拉起 -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- 后台低优先级,不抢占前台资源 -->
    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>1024</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>@HOME@/Library/Logs/neve.log</string>

    <key>StandardErrorPath</key>
    <string>@HOME@/Library/Logs/neve.error.log</string>
</dict>
</plist>
